# modules/profile/templates/testswarm/browserstack/testswarm-browserstack.service.erb
[Unit]
Description=Integration layer between TestSwarm and BrowserStack
Documentation=https://github.com/jquery/testswarm-browserstack/blob/main/README.md
After=network.target nss-lookup.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/srv/testswarm-browserstack
ExecStart=/usr/bin/node /srv/testswarm-browserstack/src/cli.js --config /etc/testswarm-browserstack.json --run-loop 60
Restart=always

PrivateTmp=full
ProtectSystem=strict
ReadOnlyPaths=/srv/testswarm-browserstack
ReadOnlyPaths=/etc/testswarm-browserstack.json

[Install]
WantedBy=multi-user.target
