# modules/notifier/templates/notifier.service.erb
[Unit]
Description=node-notifier server
Documentation=https://github.com/jquery/node-notifier-server/blob/main/README.md
After=network.target nss-lookup.target

[Service]
User=notifier
Group=notifier
WorkingDirectory=<%= @base_path %>
ExecStart=/usr/bin/node <%= @base_path %>/bin/notifier-server --debug -d /etc/notifier.d -p 3333
Restart=always

[Install]
WantedBy=multi-user.target
